# sensor-playground-web.service
[Unit]
Description=Sensor Playground Web Dashboard
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/home/pi/geepers-home
ExecStart=/bin/bash /home/pi/geepers-home/launch.sh
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

# API key for LLM chat, voice, and vision — loaded from .env by launch.sh
# Alternatively set here: Environment=DREAMER_API_KEY=your_key_here
# Get access at https://dr.eamer.dev/code/api

[Install]
WantedBy=multi-user.target
